# monitoring/quadlet/prometheus.container
# Podman Quadlet unit file for Prometheus
# Installed by monitoring-setup.sh to ~/.config/containers/systemd/prometheus.container
#
# Deployed on the CENTRAL machine only.
# Scrapes all exporters (local + tunneled remote) and evaluates alert rules.
#
# Quadlet generates a systemd unit (prometheus.service) from this file.
# See: man podman-systemd.unit(5)

[Unit]
Description=Prometheus (OSAC CI Monitoring)
After=network-online.target
Wants=network-online.target

[Container]
Image=quay.io/prometheus/prometheus:v3.4.1
ContainerName=prometheus

# Use host network so Prometheus can scrape localhost-bound exporters.
# Prometheus listens on 127.0.0.1:9091 (cockpit uses 9090).
Network=host

# Mount config and data directories.
# :Z applies the correct SELinux context for rootless Podman.
Volume=%h/.monitoring-server/config/prometheus.yml:/etc/prometheus/prometheus.yml:ro,Z
Volume=%h/.monitoring-server/config/alert-rules.yml:/etc/prometheus/alert-rules.yml:ro,Z
Volume=%h/.monitoring-server/data/prometheus:/prometheus:Z

# 30-day retention, WAL compression enabled
Exec=--config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/prometheus --storage.tsdb.retention.time=30d --web.listen-address=127.0.0.1:9091 --web.enable-lifecycle

# Run as the host user inside the container.
# keep-id maps host uid (1000) into the container's user namespace.
# User=1000 overrides the image default (nobody/65534) so the process
# can write to the data volume owned by the host user.
UserNS=keep-id
User=1000

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
